# deploy/custos.service
[Unit]
Description=Custos Endustriyel Edge Izleme Sistemi (Analytics + Dashboard + Archive)
After=network-online.target postgresql.service
Wants=network-online.target
Requires=postgresql.service

[Service]
# V11-105/K13: Type=notify + WatchdogSec — sd_notify katmani.
# READY=1 lifespan tarafindan; WATCHDOG=1 her 30 sn heartbeat task'i.
Type=notify
NotifyAccess=main
WatchdogSec=60
User=custos
WorkingDirectory=/opt/custos
EnvironmentFile=/opt/custos/.env
# H-1 (29 Nis 2026 denetim): uvicorn loopback'e bağlanır, Caddy 127.0.0.1:8000'a
# proxy yapar — TLS bypass riski kapanır. --proxy-headers + --forwarded-allow-ips
# ile X-Forwarded-For okunur (rate limit IP-bazlı sayım gerçek istemci görür).
ExecStart=/opt/custos/.venv/bin/uvicorn custos.__main__:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips 127.0.0.1
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
