# deploy/custos-critical.service
[Unit]
Description=Custos Critical Loop (Modbus Collector)
After=network-online.target postgresql.service custos.service
Wants=network-online.target
Requires=postgresql.service custos.service

[Service]
# V11-105/K13: Type=notify + WatchdogSec — sd_notify katmani.
# READY=1 main() tarafindan; WATCHDOG=1 her 30 sn heartbeat task'i.
Type=notify
NotifyAccess=main
WatchdogSec=60
User=custos
WorkingDirectory=/opt/custos
EnvironmentFile=/opt/custos/.env
ExecStart=/opt/custos/.venv/bin/python -m custos.critical
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
