# deploy/custos-critical.service
[Unit]
Description=Custos Critical Loop (Modbus Collector)
After=network-online.target postgresql.service custos.service
Wants=network-online.target
Requires=postgresql.service custos.service

[Service]
Type=simple
User=custos
WorkingDirectory=/opt/custos
EnvironmentFile=/opt/custos/.env
ExecStart=/opt/custos/.venv/bin/python -m custos.critical
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
